# just.service
[Unit]
Description = Just a Service
[Service]
ExecStart   = /opt/just.sh
[Install]
WantedBy    = multi-user.target
